# stage2/17-clean_logs/files/clean_tmp_folder.service
[Unit]
Description=Exec clean tmp folder once

[Service]
Type=simple
ExecStart=clean_tmp_logs.sh

[Install]
WantedBy=default.target
