# stage2/08-online-update-daemon/files/ota-update.service
[Unit]
Description=Pull OTA update once

[Service]
Type=simple
ExecStart=ota-update-daemon.sh
RuntimeMaxSec=3600

[Install]
WantedBy=default.target
